# contrib/python/python-crontab/tests/data/crontab

* * * * * peter parker
* 2 * * * driver parker

