# contrib/python/python-crontab/tests/data/crontabs/system_one

* * * * * bilbo baggins
* 2 * * * frodo baggins

